# apps/drive/docker/root/etc/services.d/filebrowser/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid abc devault -c /config/settings.json -d /database/devault.db;